# src/deb/upstart/schematf.conf
# vim: set ft=upstart ts=4 et:
description "schematf"

start on runlevel [2345]
stop on runlevel [!2345]

limit nofile 64000 64000

kill timeout 300 # wait 300s between SIGTERM and SIGKILL.

pre-start script
    mkdir -p /var/lib/schematf/
    mkdir -p /var/log/schematf/
end script

script
    ENABLE_SCHEMATF="yes"

    if [ "x$ENABLE_SCHEMATF" = "xyes" ]; then
        exec start-stop-daemon --start --quiet --chuid root  \
            --exec /opt/solong/schematf/current/bin/start -- 
    fi
end script
